# systemd-units/claudebox-up.service
[Unit]
Description=claudebox OWNER — start the box so its conmon lands in an INDEPENDENT scope (not a caller's cgroup)
Documentation=https://github.com/oso-gato/fedora-bootstrap
After=podman.socket network-online.target
Wants=podman.socket

[Service]
Type=oneshot
# THE FIX (incident 2026-07-11, 41 box deaths): a container started by podman FROM WITHIN a systemd
# unit does NOT get its own libpod-conmon scope — podman bails on scope creation when INVOCATION_ID is
# set (quadlet behavior). conmon then lands in the STARTING unit's cgroup, and dies with it (exit 143)
# the moment that unit's cgroup is torn down. The watcher oneshots (live-gate/host-agent), which used
# to auto-start the box via `distrobox enter`, were the killers. `env -u INVOCATION_ID` removes that
# variable for the podman child, so podman DOES create the container's own scope under the user manager
# — an independent lifetime that no watcher-tick teardown can touch.
ExecCondition=/usr/bin/bash -c '! systemctl --user is-active --quiet claudebox-rebuild-run.service'
ExecStart=/usr/bin/env -u INVOCATION_ID /usr/bin/podman start claudebox
# `podman start` is a NO-OP if the box is already running (exit 0), so this unit is safe to re-run every
# watcher tick (each watcher Wants= it) — which is the AUTOMATIC RECOVERY: if the box ever exits, the
# next tick re-runs this and revives it, always into an independent conmon scope. NOT RemainAfterExit:
# staying "active (exited)" would suppress that per-tick re-run and defeat recovery.
ExecStop=/usr/bin/podman stop -t 10 claudebox
TimeoutStartSec=120

[Install]
WantedBy=default.target
