# systemd-units/claudebox-up.service
[Unit]
Description=claudebox OWNER — start the box so its conmon lands in an INDEPENDENT scope (not a caller's cgroup)
Documentation=https://github.com/oso-gato/fedora-bootstrap
After=podman.socket network-online.target
Wants=podman.socket

[Service]
Type=oneshot
# WHY env -u INVOCATION_ID (incident 2026-07-11, 41 box deaths): a container started by podman FROM
# WITHIN a systemd unit does NOT get its own libpod-conmon scope — podman bails on scope creation when
# INVOCATION_ID is set (quadlet behavior). conmon then lands in the STARTING unit's cgroup and dies with
# it (exit 143) the moment that unit's cgroup is torn down. The 10s watcher oneshots were the killers.
# Unsetting INVOCATION_ID for the podman child makes podman create the container's OWN scope under the
# user manager — an independent lifetime no watcher-tick teardown can touch.
#
# WHY RemainAfterExit=yes AND NO ExecStop (regression fix, same day — v1.2.60 shipped BOTH mistakes):
# a Type=oneshot unit WITHOUT RemainAfterExit goes INACTIVE the instant ExecStart returns — and systemd
# runs ExecStop on that deactivation. So v1.2.60's `ExecStop=podman stop` fired ~immediately after every
# start: the unit STARTED the box and then STOPPED it, every single watcher tick (~10s). A `claude`
# session died within 5s. RemainAfterExit=yes keeps the unit "active (exited)" so no deactivation (and
# thus no stop) ever happens; carrying no ExecStop at all means the box also survives a real unit stop.
# Do NOT re-add ExecStop, and do NOT drop RemainAfterExit — either one alone re-creates the thrash.
#
# RECOVERY of a dead box is NOT this unit's job (RemainAfterExit means it will not re-fire per tick).
# It is the watchers': live-gate-watch/host-agent-watch run `env -u INVOCATION_ID … distrobox enter`
# every 10s, and `distrobox enter` STARTS a stopped box — into an independent conmon scope thanks to
# the same INVOCATION_ID belt (v1.2.60). VERIFIED LIVE on erebus 2026-07-12, not assumed: `podman stop
# claudebox` → 30s later the box was back ("Up 30 seconds") with NO human action, and it then SURVIVED
# multiple 10s watcher ticks. That survival is what proves the scope is INDEPENDENT — a conmon parented
# to a tick's cgroup dies at that tick's teardown (the incident signature: 41 deaths, each within +/-1s
# of a teardown), so outliving several ticks is only possible outside them.
RemainAfterExit=yes
ExecCondition=/usr/bin/bash -c '! systemctl --user is-active --quiet claudebox-rebuild-run.service'
ExecStart=/usr/bin/env -u INVOCATION_ID /usr/bin/podman start claudebox
TimeoutStartSec=120

[Install]
WantedBy=default.target
